# bin/gateway-daemon.service
[Unit]
Description=transmit serial port to tcpserver 62001
Requires=network-online.target
After=network.target network-online.target multi-user.target


[Service]
Type=simple
User=root
EnvironmentFile=/etc/default/gateway

ExecStart=/usr/local/bin/str2str $GW_OPTS
Restart=on-failure
RestartSec=15
KillMode=process


[Install]
WantedBy=multi-user.target
